# morse_source/feeds/morse/luci/luci-theme-morse/root/etc/uci-defaults/90_luci-theme-morse
#!/bin/sh

mv /morse/banner /etc/banner

uci set "luci.themes.Morse=/luci-static/morse"
uci set "luci.main.mediaurlbase=/luci-static/morse"
uci commit luci

exit 0
