# morse_source/feeds/morse/kernel/morse-bundle/files/etc/uci-defaults/96-morse-snapshot
#!/bin/sh 

uci set system.morse_snapshot='morse'
