# morse_source/feeds/morse/kernel/morse-bundle/files/etc/uci-defaults/98-morse-prpl-defaults
#!/bin/sh

cat > /etc/config/prplmesh <<EOF
config prplmesh 'config'
	option enable '0'
	option onboarding '1'
	option rdkb_extensions '0'
	option band_steering '1'
	option dfs_reentry '0'
	option passive_mode '1'
	option operational '1'
	option ssid 'EasyMesh'
	option mode_enabled 'WPA2-Personal'
	option wep_key '123456789a'
	option key_passphrase 'test_passphrase'
	option mem_only_psk '0'
	option backhaul_band 'auto'
	option certification_mode '0'
	option stop_on_failure_attempts '1'
	option zwdfs_enable '0'
	option best_channel_rank_th '0'
	option backhaul_wire_iface 'eth0'
	option roaming_hysteresis_percent_bonus '10'
	option steering_disassoc_timer_msec '200'
	option link_metrics_request_interval_sec '60'
	option clients_measurement_mode '1'
	option radio_stats_enable '0'
	option unsuccessful_assoc_report_policy '0'
	option unsuccessful_assoc_max_reporting_rate '30'
	option persistent_db '0'
	option clients_persistent_db_max_size '256'
	option steer_history_db_max_size '24'
	option max_timelife_delay_minutes '525600'
	option unfriendly_device_max_timelife_delay_minutes '1440'
	option persistent_db_aging_interval_sec '3600'
	option commit_changes_interval '10'
	option client_roaming '1'
	option management_mode 'Multi-AP-Agent'
	option operating_mode 'WDS-Repeater'
	option master '0'
	option gateway '0'
	option wired_backhaul '0'
EOF